# images/devshell/files/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
umask=002

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd.log

[program:jupyterlab]
command=jupyter lab --config=/home/travis/.jupyter/jupyter_server_config.py --notebook-dir=/home/travis --allow-root
user=travis
environment=HOME="/home/travis",PATH="/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/jupyterlab.log
stderr_logfile=/var/log/supervisor/jupyterlab.log

[program:codeserver]
command=code-server --bind-addr 0.0.0.0:8080 --auth none
user=travis
environment=HOME="/home/travis",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",VSCODE_PROXY_URI="./proxy/{{port}}"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/codeserver.log
stderr_logfile=/var/log/supervisor/codeserver.log

[program:web-ui]
command=python3 -m uvicorn app:app --host 0.0.0.0 --port 8090
directory=/home/travis/repos/noosphere/domains/presentation/web-ui/backend
user=travis
environment=HOME="/home/travis",PATH="/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/web-ui.log
stderr_logfile=/var/log/supervisor/web-ui.log
